# dockerfiles/centos7/drone/etc/profile.d/03-condor_environment.sh
../../../common/etc/profile.d/03-condor_environment.sh