# common/etc/profile.d/03-condor_environment.sh
#!/bin/bash

# Copyright 2022 University of Bonn
#
# SPDX-License-Identifier: GPL-3.0-or-later

echo "Hello, this is a TEST TEST TEST TEST!!!"
